# systemd/cpi.service.in
# Service unit to apply control program identification (CPI) settings
#
# Copyright IBM Corp. 2017
#
# s390-tools is free software; you can redistribute it and/or modify
# it under the terms of the MIT license. See LICENSE for details.
#

[Unit]
Description=Apply Control Program Identification (CPI)
DefaultDependencies=no
Conflicts=shutdown.target
After=sysinit.target
ConditionPathIsReadWrite=/sys/firmware/cpi

[Service]
Type=oneshot
RemainAfterExit=yes
#
# Specify a file with the environment variables using the EnvironmentFile=
# service property.
#
EnvironmentFile=@sysconf_path@/sysconfig/cpi
#
# Alternatively, you can specify the environment variables directly,
# using the Environment= service properties, which can be overridden by
# "systemctl edit --full cpi".
# In this case, the EnvironmentFile= needs to be commented out, otherwise
# it would override the directly specified variables.
#
# The sysplex name, system name, and system type consist of up to eight
# characters of the following set: A-Z, 0-9, $, @, #, and blank.
#
# Environment=CPI_SYSTEM_NAME=
# Environment=CPI_SYSPLEX_NAME=
# Environment=CPI_SYSTEM_LEVEL=
# Environment=CPI_SYSTEM_TYPE=LINUX

#
# Sending data to the HMC/SE
ExecStart=@toolslib_path@/cpictl -e

[Install]
WantedBy=multi-user.target
