# systemd/mon_fsstatd.service.in
#
# Systemd unit for mon_fsstatd (monitor z/VM guest file systems)
#
# Copyright IBM Corp. 2017
#
# s390-tools is free software; you can redistribute it and/or modify
# it under the terms of the MIT license. See LICENSE for details.
#

[Unit]
Description=Monitor z/VM guest file system size
Documentation=man:mon_fsstatd(8)
ConditionVirtualization=zvm
After=remote-fs.target

[Service]
#
# Specify a file with the environment variables using the EnvironmentFile=
# service property.
#
EnvironmentFile=@sysconf_path@/sysconfig/mon_fsstatd
#
# Alternatively, you can specify the environment variables directly,
# using the Environment= service property, which can be overridden by
# "systemctl edit --full mon_fsstatd".
# In this case, the EnvironmentFile= needs to be commented out, otherwise
# it would override the directly specified variables.
#
#Environment=FSSTAT_INTERVAL=30

ExecStartPre=-/sbin/modprobe monwriter
ExecStartPre=/usr/bin/udevadm settle --timeout=10
ExecStart=@usrsbin_path@/mon_fsstatd -i $FSSTAT_INTERVAL
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Type=forking
PIDFile=/run/mon_fsstatd.pid

[Install]
WantedBy=multi-user.target
