# etc/udev/rules.d/57-osasnmpd.rules
SUBSYSTEM=="ccwgroup", ACTION=="change", RUN+="/usr/bin/killall -SIGUSR1 osasnmpd"
